# systemd/onstep-handset.service
# REFERENCE TEMPLATE ONLY. scripts/install.sh GENERATES the real unit with paths
# and user matching your actual checkout -- prefer running install.sh. The
# User=/paths below assume user 'peter' and ~/onstepController; edit if you
# install this by hand.
[Unit]
Description=OnStep WiFi hand controller
# Ordering + requirement so the app starts once the network is actually up.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=peter
WorkingDirectory=/home/peter/onstepController
# Full paths required. Uses the project venv created by scripts/install.sh.
ExecStart=/home/peter/onstepController/.venv/bin/python -m onstep_handset.main
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
